# hypr/services/pscircle-wallpaper.timer
[Unit]
Description=Run pscircle wallpaper updates periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=1s

[Install]
WantedBy=timers.target